# ops/debian/systemd/anlagen-logrotate.service.template
[Unit]
Description=Anlagenserver log and trace rotation

[Service]
Type=oneshot
ExecStart=/bin/bash __INSTALL_ROOT__/ops/scripts/rotate_logs.sh
