# ops/debian/systemd/anlagen-logrotate.timer
[Unit]
Description=Run Anlagenserver log rotation daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
